# install/upstart.conf
# simple uWSGI script

description "uwsgi tiny instance"
start on runlevel [2345]
stop on runlevel [06]

exec /opt/anaf/env/bin/uwsgi --ini /opt/anaf/uwsgi.ini --die-on-term >> /var/log/anaf.log 2>&1